# rootfs/etc/services.d/apache/run
#!/usr/bin/with-contenv bash
if [ ! -d /var/log/aegir ]; then
    mkdir /var/log/aegir
fi

if [ ! -d /var/log/config ]; then
    mkdir -p /var/aegir/config
fi

if [ ! -d /var/aegir/config/static ]; then
    mkdir -p /var/aegir/config/static
fi

if [ ! -f /var/aegir/config/apache.conf ]; then
    echo "" > /var/aegir/config/apache.conf
fi

if [ -f /var/aegir/config/apache.conf ]; then
    ln -sf /var/aegir/config/apache.conf /etc/apache2/conf-available/aegir.conf
    ln -sf /etc/apache2/conf-available/aegir.conf /etc/apache2/conf-enabled/aegir.conf
fi
# inpired by https://github.com/wodby/apache/blob/master/docker-entrypoint.sh
_gotpl() {
    if [[ -f "$1" ]]; then
        echo "generate config from $1 to $2" 
        gotpl "$1" > "$2"
    fi
}
# apache conf
_gotpl "/etc/gotpl/apache2-mpm-prefork.tmpl" "/etc/apache2/conf-enabled/mpm-prefork.conf"

echo "starting apache..."
exec apachectl -D FOREGROUND
